# systemd/nycu-platform.service
[Unit]
Description=NYCU Course Platform - Docker Compose Services
Documentation=https://nymu.com.tw/docs
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/nycu-platform

# Service startup
ExecStartPre=/usr/bin/docker-compose -f docker-compose.yml pull
ExecStart=/usr/local/bin/docker-compose -f docker-compose.yml up

# Service shutdown
ExecStop=/usr/local/bin/docker-compose -f docker-compose.yml down

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitInterval=60s
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nycu-platform

# Resource limits
MemoryLimit=16G
CPUQuota=3200%

[Install]
WantedBy=multi-user.target
